# etc/anacrontab
# /etc/anacrontab: configuration file for anacron

# See anacron(8) and anacrontab(5) for details.

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
HOME=/root
LOGNAME=root

# These replace cron's entries
1	0	cron.daily	run-parts --report /etc/cron.daily
7	0	cron.weekly	run-parts --report /etc/cron.weekly
@monthly	0	cron.monthly	run-parts --report /etc/cron.monthly

# Backup-Manager
1	0	backup-manager	/backup-manager >/proc/1/fd/1 2>/proc/1/fd/2
