# etc/s6-overlay/s6-rc.d/crond/run
#!/command/with-contenv sh

mkdir -p /etc/cron.d
echo "${BM_CRON} root /command/backup-manager" >> /etc/cron.d/backup-manager

cron -f -l 5
